# supervisord.conf
[supervisord]
nodaemon=true

[program:celery-batch_update]
command=celery -A src.tasks.batch.update.tasks worker --loglevel=info -Q batch_update --concurrency=10 --statedb=/app/worker-batch_update.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-batch_update.log
stdout_logfile=/app/app-batch_update.log

[program:celery-homelette]
command=celery -A src.tasks.homelette.tasks worker --loglevel=info -Q homelette --concurrency=10 --statedb=/app/worker-homelette.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-homelette.log
stdout_logfile=/app/app-homelette.log

[program:celery-musite]
command=celery -A src.tasks.musite.tasks worker --loglevel=info -Q musite --concurrency=10 --statedb=/app/worker-musite.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-musite.log
stdout_logfile=/app/app-musite.log

[program:celery-boltz]
command=celery -A src.tasks.boltz.tasks worker --loglevel=info -Q boltz --concurrency=10 --statedb=/app/worker-boltz.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-boltz.log
stdout_logfile=/app/app-boltz.log

[program:celery-diffab]
command=celery -A src.tasks.diffab.tasks worker --loglevel=info -Q diffab --concurrency=10 --statedb=/app/worker-diffab.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-diffab.log
stdout_logfile=/app/app-diffab.log

[program:celery-gnina]
command=celery -A src.tasks.gnina.tasks worker --loglevel=info -Q gnina --concurrency=10 --statedb=/app/worker-gnina.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-gnina.log
stdout_logfile=/app/app-gnina.log

[program:celery-gan]
command=celery -A src.tasks.gan worker --loglevel=info -Q gan --concurrency=10 --statedb=/app/worker-gan.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-gan.log
stdout_logfile=/app/app-gan.log

[program:celery-blood-test]
command=celery -A src.tasks.blood_test worker --loglevel=info -Q blood_test --concurrency=10 --statedb=/app/worker-blood-test.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-blood-test.log
stdout_logfile=/app/app-blood-test.log

[program:celery-url-download]
command=celery -A src.tasks.url_download worker --loglevel=info -Q url_download --concurrency=10 --statedb=/app/worker-url-download.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-url-download.log
stdout_logfile=/app/app-url-download.log

[program:celery-vep-ensemble]
command=celery -A src.tasks.vep_ensemple worker --loglevel=info -Q vep_ensemble --concurrency=10 --statedb=/app/worker-vep-ensemble.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-vep-ensemble.log
stdout_logfile=/app/app-vep-ensemble.log

[program:celery-variant-calling-bam]
command=celery -A src.tasks.variant_calling_bam worker --loglevel=info -Q bam_processing --concurrency=10 --statedb=/app/worker-variant-calling-bam.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-variant-calling-bam.log
stdout_logfile=/app/app-variant-calling-bam.log

[program:celery-variant-calling-fastq]
command=celery -A src.tasks.variant_calling_fastq worker --loglevel=info -Q fastq_processing --concurrency=10 --statedb=/app/worker-variant-calling-fastq.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-variant-calling-fastq.log
stdout_logfile=/app/app-variant-calling-fastq.log

[program:celery-microbiome]
command=celery -A src.tasks.microbiome worker --loglevel=info -Q microbiome --concurrency=10 --statedb=/app/worker-microbiome.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-microbiome.log
stdout_logfile=/app/app-microbiome.log

[program:celery-genetic-annotation]
command=celery -A src.tasks.genetic_annotation worker --loglevel=info -Q genetic_annotation --concurrency=10 --statedb=/app/worker-genetic-annotation.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-genetic-annotation.log
stdout_logfile=/app/app-genetic-annotation.log

[program:celery-genetic]
command=celery -A src.tasks.genetic worker --loglevel=info -Q genetic --concurrency=10 --statedb=/app/worker-genetic.state
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/app-genetic.log
stdout_logfile=/app/app-genetic.log
